# roles/linux_security/templates/apt-daily-upgrade.timer.j2
# {{ ansible_managed }}
[Timer]
OnCalendar=
OnCalendar=*-*-* {{ linux_security_apt_upgrade_time }}
RandomizedDelaySec={{ linux_security_apt_upgrade_randomized_delay_sec }}
