# roles/linux_security/templates/apt-daily.timer.j2
# {{ ansible_managed }}
[Timer]
OnCalendar=
OnCalendar=*-*-* {{ linux_security_apt_download_time }}
RandomizedDelaySec={{ linux_security_apt_download_randomized_delay_sec }}
